# services/claude-hub.service
[Unit]
Description=claude-hub — local landing page + path-routed reverse proxy
After=network.target

[Service]
Type=simple
# EDIT BEFORE INSTALLING: User= must match the account that owns the repo,
# and the bash login shell must source nvm so `node` resolves to your
# preferred version. marked v18 is ESM-only → node ≥ 22 required.
User=david
# systemd does NOT expand %h to User=david's home in system-mode units —
# %h resolves against the manager (root). Use /home/<user> explicitly, or
# replace this path with your account's home before installing.
WorkingDirectory=/home/USER/projects/claude-hub
ExecStart=/bin/bash -lc 'export NVM_DIR=$HOME/.nvm && . $NVM_DIR/nvm.sh && exec node server.js'
# KillMode=process: stop only the main node PID on restart. Child ttyd
# subprocesses (and the tmux servers they forked) escape the cgroup and
# keep running, so a `systemctl restart claude-hub.service` does NOT kill
# live Claude tmux sessions. Server.js unlinks stale sockets on startup so
# the next /term/<name>/ hit can rebind. Orphan ttyd procs are unreachable
# (their socket path is gone) and live until host reboot — small leak,
# trade-off for not nuking user sessions.
KillMode=process
Restart=always
RestartSec=2
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
